# supervisord.conf
[supervisord]
nodaemon=true

[program:openvpn]
command=openvpn

[program:tinyproxy]
command=tinyproxy -c /etc/openvpn/tinyproxy.conf
